# contrib/debian/DEBIAN.agent/prerm
#!/bin/sh










# remove service init startup
update-rc.d -f statcollector-agent remove

# shutdown service
service statcollector-agent stop

# EOF